# misc/profile.d/mapmo.sh
redshift -O 5000
ksuperkey
plank &>/dev/null &
